# etc/startup/alt-tab.conf
start on sway-start
stop on sway-end
respawn
exec alt-tab-daemon
