# etc/startup/alt-tab.conf
start on sway-session
stop on stopping sway
respawn
exec alt-tab-daemon
