# etc/startup/mako.conf
start on sway-start and started session-bus
stop on sway-end or stopping session-bus
respawn
reload signal 0
exec mako
