# etc/startup/mako.conf
start on started sway and started session-bus
stop on stopping sway or stopping session-bus
respawn
reload signal 0
exec mako
