# etc/startup/screen-dpms.conf
start on swayidle event=timeout timeout=screen
stop on swayidle event=resume resume=screen
task
pre-start exec swaymsg output '*' dpms off
post-stop exec swaymsg output '*' dpms on
