# etc/startup/screen-off.conf
start on sway-long-idle
task
exec swaymsg output '*' dpms off
