# etc/startup/screen-on.conf
start on sway-active
task
exec swaymsg output '*' dpms on
