# etc/startup/screen-dpms.conf
start on sway-long-idle
stop on sway-active
task
pre-start swaymsg output '*' dpms off
pre-stop swaymsg output '*' dpms on
