# etc/startup/smooth-fading.conf
start on sway-session
stop on stopping sway
respawn
exec smooth-fading-sway.py
