# etc/startup/snooze-newsboat.conf
start on startup
stop on end-session
respawn
exec snooze-newsboat
