# etc/startup/snooze-newsboat.conf
start on session-start
stop on session-end
respawn
exec env -u JOB \
		 -u UPSTART_JOB \
		 -u UPSTART_EVENTS \
		 -u UPSTART_INSTANCE \
		 -u UPSTART_DIR \
		 snooze-newsboat
