# etc/startup/snooze-newsboat.conf
start on startup
stop on session-end
respawn
exec snooze-newsboat
