# etc/startup/telegram-desktop.conf
start on sway-start
stop on sway-end
respawn
env QT_QPA_PLATFORM=wayland
exec env -u JOB \
		 -u UPSTART_JOB \
		 -u UPSTART_EVENTS \
		 -u UPSTART_INSTANCE \
		 -u UPSTART_DIR \
		 telegram-desktop
