# etc/startup/telegram-desktop.conf
start on sway-session
stop on stopping sway
respawn
env QT_QPA_PLATFORM=wayland
exec telegram-desktop
