# etc/startup/telegram-desktop.conf
start on started sway
stop on stopping sway
respawn
env QT_QPA_PLATFORM=wayland
exec telegram-desktop
