# etc/startup/waybar.conf
start on sway-start
stop on sway-end
respawn
exec env -u JOB \
		 -u UPSTART_JOB \
		 -u UPSTART_EVENTS \
		 -u UPSTART_INSTANCE \
		 -u UPSTART_DIR \
		 waybar
