# etc/startup/waybar.conf
start on started sway
stop on stopping sway
respawn
env LANG=C
exec env -u JOB \
		 -u UPSTART_JOB \
		 -u UPSTART_EVENTS \
		 -u UPSTART_INSTANCE \
		 -u UPSTART_DIR \
		 waybar
